# root/etc/services.d/ldap-daemon/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid python /app/nginx-ldap-auth-daemon.py
